# deploy/systemd/simrp-api.service.example
[Unit]
Description=SIMRP API prototype service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/simrp
Environment=SIMRP_ENV=production
Environment=SIMRP_HOST=127.0.0.1
Environment=SIMRP_PORT=8000
Environment=SIMRP_ENABLE_DEMO_SEED=false
Environment=SIMRP_ALLOWED_ORIGINS=https://simrp.example.go.id
EnvironmentFile=/etc/simrp/simrp.env
ExecStart=/usr/bin/python3 /opt/simrp/server/main.py
Restart=on-failure
RestartSec=5
User=simrp
Group=simrp
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/simrp/database/runtime /opt/simrp/database/backups

[Install]
WantedBy=multi-user.target
